# tb3-mcp/deploy/tb3-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.tb3.mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>/ABSOLUTE/PATH/TO/TB3-ESP32/tb3-mcp/dist/server.js</string>
  </array>
  <key>WorkingDirectory</key><string>/ABSOLUTE/PATH/TO/TB3-ESP32/tb3-mcp</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/tb3-mcp.log</string>
  <key>StandardErrorPath</key><string>/tmp/tb3-mcp.err</string>
</dict>
</plist>
